# installer/data/forge-tui.service
[Unit]
Description=InterGenOS Forge Installer — Declarative Builder TUI
Documentation=https://github.com/InterGenJLU/intergenos/blob/main/installer/README.md
After=multi-user.target
Wants=multi-user.target
ConditionKernelCommandLine=igos.mode=install-tui

# Claim /dev/tty1 exclusively. Without these, getty@tty1.service from
# the default getty.target.wants/ set would race us for tty1; agetty
# typically wins, leaving forge-tui without a usable terminal.
# Conflicts= stops the getty unit when we start; Before= sequences us
# to start first so the getty never gets the chance.
Conflicts=getty@tty1.service
Before=getty@tty1.service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/forge --mode tui --archives /var/lib/igos/archives --packages /var/lib/igos/packages
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TimeoutStopSec=10
Restart=no
IgnoreSIGPIPE=no
KillMode=process
SendSIGHUP=yes

[Install]
WantedBy=multi-user.target
