# installer/data/intergenos-first-boot-greeter.service
[Unit]
Description=InterGenOS First-Boot Password Greeter
Documentation=man:systemd.unit(5)

# Wait for vconsole + tmpfiles setup before claiming tty1 for the
# password prompts. Conflicts with getty@tty1.service locks the TTY.
# Stale ordering on intergen-firstboot.service removed; the firstboot
# rewrite uses a per-user autostart .desktop chain on first GDM login,
# not a systemd unit.
After=systemd-vconsole-setup.service
After=systemd-tmpfiles-setup.service

# Idempotency: do nothing on subsequent boots once the flag is written.
ConditionPathExists=!/etc/intergenos/first-boot-completed

# Live-installer mode guard. The squashfs is shared between live-boot and
# installed-system rootfs, so this unit ships enabled in both contexts. On
# live-ISO boot the kernel cmdline carries `igos.mode=live`, `=install-gui`
# or `=install-tui` (per installer/init/cmdline.*.txt); the greeter must NOT
# fire there — those flows claim tty1 themselves (forge-tui.service on tty1
# for install-tui, GDM/forge-gui session on graphical.target for live and
# install-gui). Negated bare-word match: condition is TRUE only when the
# cmdline has no `igos.mode` assignment, which is the installed-system case.
ConditionKernelCommandLine=!igos.mode

DefaultDependencies=no
Conflicts=getty@tty1.service
Before=getty@tty1.service

# If the greeter ever fails, fall through to getty rather than blocking
# boot indefinitely. Path 4 ensured the image has builder-set passwords
# as the safety net — getty is reachable with those if the greeter is
# bypassed (next boot will retry the greeter).
OnFailure=getty@tty1.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/intergenos/first-boot-greeter
RemainAfterExit=no

# Block the TTY for our prompts.
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

# No timeout — user may take their time choosing passwords.
TimeoutStartSec=0

# Medium hardening (19 directives) per audit rows F-048 + F-049 surgical-
# hardening exemption closure. The greeter runs as oneshot root from
# systemd; password-tool needs CAP_CHOWN + CAP_DAC_OVERRIDE + CAP_FOWNER
# to write the credential database, and CAP_AUDIT_WRITE for PAM audit
# events via libaudit (which opens AF_NETLINK sockets — both AF_UNIX and
# AF_NETLINK must remain in RestrictAddressFamilies for the audit trail
# to land; without AF_NETLINK, libaudit silently degrades and CAP_AUDIT_
# WRITE becomes inert). ProtectSystem=true keeps /usr + /boot read-only
# but leaves /etc writable (greeter writes credential database + flag at
# /etc/intergenos/first-boot-completed). SystemCallFilter=@system-service
# is the systemd.exec(5) predefined syscall group for typical service
# workloads (bash + native binary + no JIT engines).
NoNewPrivileges=yes
ProtectSystem=true
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHostname=yes
PrivateTmp=yes
RestrictAddressFamilies=AF_UNIX AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
CapabilityBoundingSet=CAP_CHOWN CAP_DAC_OVERRIDE CAP_FOWNER CAP_AUDIT_WRITE
MemoryDenyWriteExecute=yes

[Install]
WantedBy=multi-user.target
